# anbox-launchers.service
[Unit]
Description=Anbox Desktop Entry Categories Fixer.

[Service]
Type=oneshot
ExecStart=@PREFIX@/anbox/anbox-launchers
